# project/chemtools.conf
[program:chemtools]
command=/home/chris/project/bin/gunicorn project.wsgi:application
directory=/home/chris/project/chemtools-webapp
user=chris
autostart=true
autorestart=true
redirect_stderr=true
