# project/supervisor_settings.conf
[program:chemtools]
command=$CHEMTOOLS_DIR/bin/gunicorn project.wsgi:application
directory=$CHEMTOOLS_DIR
user=$INSTALL_USER
autostart=true
autorestart=true
redirect_stderr=true